# twitter_bot_lastfm_cg/systemd-service/twitter_bot_lastfm_cg_weekly.timer
[Unit]
Description=lastfm_cg_weekly.timer

[Timer]
OnCalendar=*-*-* 9:00:00
Persistent=true

[Install]
WantedBy=timers.target
