# systemd-service/lastfm_cg_semiannually.timer
[Unit]
Description=lastfm_cg_semiannually.timer

[Timer]
OnCalendar=semiannually
Persistent=true

[Install]
WantedBy=timers.target
