# systemd-service/lastfm_cg.timer
[Unit]
Description=lastfm_cg.timer

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
