# systemd-service/lastfm_cg_weekly.timer
[Unit]
Description=lastfm_cg_weekly.timer

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
