# systemd-service/lastfm_cg_yearly.timer
[Unit]
Description=lastfm_cg_yearly.timer

[Timer]
OnCalendar=yearly
Persistent=true

[Install]
WantedBy=timers.target
