# bot_lastfm_cg/systemd-service/bot_lastfm_cg_monthly.timer
[Unit]
Description=bot_lastfm_cg_monthly.timer

[Timer]
OnCalendar=*-*-1 08:30:00
Persistent=true

[Install]
WantedBy=timers.target
