# bot_lastfm_cg/systemd-service/bot_lastfm_cg_yearly.timer
[Unit]
Description=bot_lastfm_cg_yearly.timer

[Timer]
OnCalendar=*-1-1 08:00:00
Persistent=true

[Install]
WantedBy=timers.target
